# tf/files/builder.conf
# builder.conf
start on networking
exec /home/ubuntu/builder server
